# crontab
# START CRON JOB
0 */12 * * * cd /Mfinanzas_app && /usr/bin/python3.11 scr/main.py >> cron_log.log 2>&1
# END CRON JOB